# deploy/systemd/discoverd.service
[Unit]
Description=discoverd UDP multicast discovery server
Documentation=https://github.com/napilab/discoverd
Wants=network-online.target
After=network-online.target

[Service]
Type=simple
User=discoverd
Group=discoverd
Environment=DISCOVERD_MCAST=239.255.255.250
Environment=DISCOVERD_PORT=9999
Environment=DISCOVERD_MAX_CLOCK_SKEW=10s
EnvironmentFile=/etc/discoverd
ExecStartPre=/usr/bin/test -n "${DISCOVERD_SECRET}"
ExecStart=discoverd --mode server --mcast ${DISCOVERD_MCAST} --port ${DISCOVERD_PORT} --max-clock-skew ${DISCOVERD_MAX_CLOCK_SKEW}
Restart=on-failure
RestartSec=2s
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true
LockPersonality=true
MemoryDenyWriteExecute=true
RestrictAddressFamilies=AF_UNIX AF_INET AF_INET6
CapabilityBoundingSet=
AmbientCapabilities=
UMask=0077

[Install]
WantedBy=multi-user.target
